# .ansible/roles/web-services/templates/syncthing.service.j2
# NOTE: THIS FILE WAS AUTO-GENERATED BY ANSIBLE
# Edit the template file in .ansible/roles/web-services/templates/syncthing.service.j2

[Unit]
Description=Syncthing - Open Source Continuous File Synchronization
After=network.target

[Service]
User={{ user_name }}
ExecStart=/usr/bin/syncthing -no-browser -no-restart -logflags=0
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
