# .ansible/roles/python/templates/jupyterhub.service.j2
# NOTE: THIS FILE WAS AUTO-GENERATED BY ANSIBLE
# Edit the template file in .ansible/roles/python/templates/jupyterhub.service.j2

[Unit]
Description=JupyterHub
After=network.target

[Service]
Type=simple
User={{ user_name }}
ExecStart=/usr/bin/jupyterhub
WorkingDirectory={{ home_dir }}
Restart=on-failure

[Install]
WantedBy=multi-user.target
