# .ansible/roles/terminal/templates/pacmanquery.service.j2
# NOTE: THIS FILE WAS AUTO-GENERATED BY ANSIBLE
# Edit the template file in .ansible/roles/terminal/templates/pacmanquery.service.j2

[Unit]
Description=Check for package updates

[Service]
Type=oneshot
ExecStart=/usr/bin/pacman -Qu

[Install]
WantedBy=multi-user.target
