# .ansible/roles/terminal/templates/tmux.service.j2
# NOTE: THIS FILE WAS AUTO-GENERATED BY ANSIBLE
# Edit the template file in .ansible/roles/terminal/templates/tmux.service.j2

[Unit]
Description=Start tmux in detached session

[Service]
Type=forking
User={{ user_name }}
ExecStart=/usr/bin/tmux new-session -s %u -d
ExecStop=/usr/bin/tmux kill-session -t %u

[Install]
WantedBy=multi-user.target
