# .ansible/roles/networking/templates/openvpn.service.j2
# NOTE: THIS FILE WAS AUTO-GENERATED BY ANSIBLE
# Edit the template file in .ansible/roles/networking/templates/openvpn.service.j2

[Unit]
Description=OpenVPN connection
After=network.target

[Service]
Type=simple
User={{ user_name }}
ExecStart=/usr/bin/openvpn --config {{ home_dir }}/.config/openvpn/%i.ovpn
Restart=on-failure

[Install]
WantedBy=multi-user.target
